# packaging/systemd/monitor-privileged-snapshot.timer
[Unit]
Description=Refresh privileged monitor snapshot every 2 minutes

[Timer]
OnUnitActiveSec=2min
AccuracySec=15s
Unit=monitor-privileged-snapshot.service

[Install]
WantedBy=timers.target
